# cookbooks/deckbrew/templates/default/deckbrew-api.conf.erb
description "DeckBrew API"

respawn

start on <%= node['deckbrew']['event'] %>
env DECKBREW_HOSTNAME=<%= node['deckbrew']['hostname'] %>
env DATABASE_USER="<%= node['deckbrew']['database']['user'] %>"
env DATABASE_HOST="<%= node['deckbrew']['database']['host'] %>"
env DATABASE_PASSWORD="<%= node['deckbrew']['database']['password'] %>"

script
        chdir /usr/local/deckbrew
        exec ./brewapi
end script

